# home/chris/.config/systemd/user/gammastep-indicator.service.d/override.conf
[Unit]
After=sway-session.target

[Service]
ExecStart=
ExecStart=/usr/bin/gammastep-indicator -l 52:13

[Install]
WantedBy=sway-session.target