# etc/systemd/user/mpv@.service
[Unit]
OnFailure=notify@%n

[Service]
ExecStart=/usr/bin/mpv --playlist=%I --shuffle --loop-playlist --pause
