# home/chris/.config/systemd/user/git-pull@.timer
[Unit]
Description=pull git changes

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=default.target
