# home/chris/.config/systemd/user/notify-failure@.service
[Service]
ExecStart=/usr/bin/notify-send -u critical -t 7000 failed "%I"
