# home/chris/.config/systemd/user/picom.service
[Unit]
Description=Picom
OnFailure=notify@%n

[Service]
ExecStart=/usr/bin/picom --experimental-backends --config=${HOME}/.config/picom.conf
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
