# home/chris/.config/systemd/user/polybar.service
[Unit]
Description=An information bar
OnFailure=notify@%n

[Service]
ExecStart=/usr/bin/polybar schwarzerblock
