# home/chris/.config/systemd/user/random-background@.path
[Unit]
Description=Sets random sway backgrounds

[Path]
PathChanged=%I
Unit=random-background@.timer

[Install]
WantedBy=sway-session.target
