# home/chris/.config/systemd/user/random-background@.timer
[Unit]
Description=Sets random sway backgrounds

[Timer]
OnCalendar=hourly
Persistent=false

[Install]
WantedBy=sway-session.target
