# home/chris/.config/systemd/user/update-mpd.timer
[Unit]
Description=mpd database update

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=default.target
