# home/chris/.config/systemd/user/wallpaper.service
[Unit]
Description=Wallpaper
Documentation=man:feh(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/feh --randomize --bg-fill /data/sync/wallpapers
